# scripts/systemd/premiumhub-apis-order-injector.service
[Unit]
Description=PremiumHub Auto Dummy FiveSim Order Injector (apis@lokal.com)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/premium-hub
Environment=APP_ROOT=/home/ubuntu/premium-hub/premiumhub-api
Environment=ENV_FILE=/home/ubuntu/premium-hub/premiumhub-api/.env
Environment=TARGET_EMAIL=apis@lokal.com
Environment=PROVIDER_STATUS=TIMEOUT
Environment=ORDER_TYPE=activation
Environment=ROWS_PER_TICK=1
Environment=LOOP_INTERVAL_SEC=1
Environment=MAX_INSERT=0
Environment=ENABLE_RESOLVED_AT=1
Environment=PGLOCK_TIMEOUT_MS=250
Environment=PGSTATEMENT_TIMEOUT_MS=2500
EnvironmentFile=-/etc/default/premiumhub-apis-order-injector
ExecStart=/home/ubuntu/premium-hub/scripts/ops/auto-inject-fivesim-orders.sh
Restart=always
RestartSec=1
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
