# scripts/systemd/workspace-premiumhub-api.service
[Unit]
Description=workspace-PremiumHub API (Go)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/premium-hub-workspace/premium-hub/premiumhub-api
EnvironmentFile=/home/ubuntu/.config/workspace-premiumhub-api.env
Environment=GIN_MODE=release
ExecStart=/home/ubuntu/premium-hub-workspace/premium-hub/premiumhub-api/bin/workspace-premiumhub-api
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
